# deploy/systemd/cloudflared.service
[Unit]
Description=Punaro Cloudflare Tunnel connector
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=cloudflared
Group=cloudflared
LoadCredential=tunnel-token:/etc/punaro/credentials/cloudflared-tunnel-token
ExecStart=/usr/local/libexec/punaro/run-cloudflared
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
ProtectProc=invisible
ProcSubset=pid
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=yes
SystemCallArchitectures=native
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077

[Install]
WantedBy=multi-user.target
