# deploy/systemd/punarod.service
[Unit]
Description=Punaro agent relay
After=network-online.target
Wants=network-online.target
[Service]
Type=exec
User=punaro
Group=punaro
StateDirectory=punaro
Environment=PUNARO_LISTEN_ADDR=127.0.0.1:8080
Environment=PUNARO_DATA_DIR=/var/lib/punaro
EnvironmentFile=-/etc/punaro/punaro.env
ExecStart=/usr/local/bin/punarod
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/var/lib/punaro
UMask=0077
[Install]
WantedBy=multi-user.target
