# deploy/systemd/punaro-telegram.service
[Unit]
Description=Punaro Telegram gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=punaro-telegram
Group=punaro-telegram
StateDirectory=punaro-telegram
EnvironmentFile=/etc/punaro/telegram.env
LoadCredential=telegram-bot-token:/etc/punaro/credentials/telegram-bot-token
LoadCredential=telegram-access:/etc/punaro/credentials/telegram-access
Environment=PUNARO_TELEGRAM_BOT_TOKEN_FILE=%d/telegram-bot-token
Environment=PUNARO_TELEGRAM_ACCESS_TOKEN_FILE=%d/telegram-access
ExecStart=/usr/local/bin/punaro-telegram
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/var/lib/punaro-telegram
UMask=0077

[Install]
WantedBy=multi-user.target
